# packages/xmedia/dvb/vtuner-apps/system.d/vtunerc.service
[Unit]
Description=Vtuner Client
Requires=network-online.target
After=network-online.target
Before=vdr.service
ConditionPathExists=/storage/.cache/services/vtunerc.conf

[Service]
ExecStart=/bin/sh -c 'exec /usr/bin/vtunerc.start'
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/sbin/rmmod -f vtunerc
TimeoutStopSec=1s
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
