# packages/xmedia/libretro/retroarch/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory RetroArch
Before=retroarch.service
ConditionPathExists=/var/lock/start.games

[Mount]
What=none
Where=/tmp/joypads
Type=overlay
Options=lowerdir=/etc/retroarch-joypad-autoconfig,upperdir=/storage/retroarch/joypads,workdir=/storage/.tmp/joypads-workdir

[Install]
WantedBy=kodi.target
